# deploy/node_b.service.in
[Unit]
Description=Accelerometer Node B
After=network.target

[Service]
Type=simple
ExecStart=@NODE_B_BINARY_PATH@ @SERVER_PORT@
WorkingDirectory=@PROJECT_ROOT_PATH@
Restart=on-failure
TimeoutStopSec=5

StandardOutput=journal
StandardError=journal
Environment=LD_LIBRARY_PATH=@LD_LIBRARY_PATH@

[Install]
WantedBy=default.target
