# deploy/server.service.in
[Unit]
Description=Accelerometer gRPC Server
After=network.target

[Service]
Type=simple
ExecStart=@SERVER_BINARY_PATH@ @SERVER_PORT@
WorkingDirectory=@PROJECT_ROOT_PATH@
Restart=on-failure

StandardOutput=journal
StandardError=journal
Environment=LD_LIBRARY_PATH=@LD_LIBRARY_PATH@

[Install]
WantedBy=default.target
